# src/_installScripts/_Configs/MotdDefaultConfig-conf/etc/cron.daily/motd
#!/bin/bash
###############################################################################
# DO NOT EDIT, THIS FILE CAN BE UPDATED WITHOUT NOTICE
###############################################################################

update-motd
