# conf/etc/update-motd.d/00-wsl-header
#!/bin/bash

echo "$(tput setaf 2)   .~~.   .~~.    $(date +"%A, %e %B %Y, %r")"
echo "$(tput setaf 2)  '. \ ' ' / .'   $(uname -srmo)"
