# src/_installScripts/_Configs/ShellBashDefaultConfig-conf/.bash-dev-env/profile.d/fixWslDate.sh
#!/bin/bash

sudo "${HOME}/.bash-dev-env/fixWslDate"
