# src/_installScripts/Aws/Awsume-conf/.bash-dev-env/profile.d/awsume.zsh
#!/usr/bin/env zsh
###############################################################################
# DO NOT EDIT, THIS FILE CAN BE UPDATED WITHOUT NOTICE
###############################################################################

#Auto-Complete function for AWSume
fpath=(~/.awsume/zsh-autocomplete/ $fpath)
