# src/_binaries/BashTools/conf/.bash-dev-env/profile.d/bash-tools.sh
#!/bin/bash
###############################################################################
# DO NOT EDIT, THIS FILE CAN BE UPDATED WITHOUT NOTICE
###############################################################################

[[ :${PATH}: == *":${HOME}/fchastanet/bash-tools/bin:"* ]] ||
  PATH="${HOME}/fchastanet/bash-tools/bin:${PATH}"

export PATH
