# src/_binaries/Docker/conf/.bash-dev-env/profile.d/docker.sh
#!/bin/bash
###############################################################################
# DO NOT EDIT, THIS FILE CAN BE UPDATED WITHOUT NOTICE
###############################################################################

if command -v docker &>/dev/null; then
  export DOCKER_HOST='unix:///run/docker.sock'
fi
