# src/_installScripts/Aws/AwsCli-conf/.bash-dev-env/profile.d/saml2aws.sh
#!/usr/bin/env bash
###############################################################################
# DO NOT EDIT, THIS FILE CAN BE UPDATED WITHOUT NOTICE
###############################################################################

export DBUS_SESSION_BUS_ADDRESS=/dev/null
