# conf/etc/cron.d/bash-dev-env-upgrade
# /etc/cron.d/upgrade: crontab entries for the bash-dev-env upgrade

SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/games:/usr/local/games:/usr/lib/wsl/lib:/mnt/c/Windows/system32:/mnt/c/Windows:/mnt/c/Windows/System32/Wbem:/mnt/c/Windows/System32/WindowsPowerShell/v1.0/

# Upgrade every 7 days at 22pm
00 22 */7 * * root  @COMMAND@
