# src/_binaries/MotdDefaultConfig/conf/etc/update-motd.d/00-wsl-header
#!/bin/bash
###############################################################################
# DO NOT EDIT, THIS FILE CAN BE UPDATED WITHOUT NOTICE
###############################################################################

echo "$(tput setaf 2)   .~~.   .~~.    $(date +"%A, %e %B %Y, %r")"
echo "$(tput setaf 2)  '. \ ' ' / .'   $(uname -srmo)"
