# src/_binaries/MotdDefaultConfig/conf/etc/update-motd.d/03-wsl-automatic-upgrade
#!/bin/bash
###############################################################################
# DO NOT EDIT, THIS FILE CAN BE UPDATED WITHOUT NOTICE
###############################################################################

BASH_DEV_ENV_LOGS_DIR=/invalid

# pam_motd does not carry the environment
if [[ -f /etc/default/locale ]]; then
  source /etc/default/locale
fi
export LANG

echo
if [[ -f "${BASH_DEV_ENV_LOGS_DIR}/lastInstallStatus" ]]; then
  cat "${BASH_DEV_ENV_LOGS_DIR}/lastInstallStatus"
else
  echo "$(tput setaf 1)(*) Weekly Automatic upgrade does not seems to run$(tput sgr0)"
fi
