# conf/Composer/etc/profile.d/composer_path.sh
#!/bin/bash

[[ :${PATH}: == *":/home/wsl/.config/composer/vendor/bin:"* ]] || PATH="/home/wsl/.config/composer/vendor/bin:${PATH}"
export PATH